# policyd/trunk/src/contrib/policyd.cron
# clean up policyd database
0 * * * * /usr/local/policyd/cleanup -c /usr/local/policyd/policyd.conf
